# systemd/backup-to-pve.timer
[Unit]
Description=Disparo diario del backup vigilancia -> pve

[Timer]
OnCalendar=*-*-* 02:30:00
# OnCalendar es un EJEMPLO: cambia la hora en tu despliegue (y no la publiques).
# Si el host estaba apagado a esa hora (apagon/NTP), correr al volver. RES-9.
Persistent=true
RandomizedDelaySec=120
Unit=backup-to-pve.service

[Install]
WantedBy=timers.target
