# systemd/check-backup-pve.timer
[Unit]
Description=Chequeo diario del backup vigilancia -> pve (heartbeat Telegram)

[Timer]
OnCalendar=*-*-* 07:00:00
# OnCalendar es un EJEMPLO: cambia la hora en tu despliegue.
# Si el host estaba apagado a esa hora, correr al volver.
Persistent=true
RandomizedDelaySec=60
Unit=check-backup-pve.service

[Install]
WantedBy=timers.target
