# systemd/frigate-weekly-restart.timer
[Unit]
Description=Restart semanal de Frigate

[Timer]
OnCalendar=Sun *-*-* 04:00:00
# OnCalendar es un EJEMPLO: cambia dia/hora en tu despliegue.
# Si el LXC estaba apagado a esa hora, correr al volver. RES-9.
Persistent=true
RandomizedDelaySec=300
Unit=frigate-weekly-restart.service

[Install]
WantedBy=timers.target
